# deploy/systemd/webauthnp256-alert.service
# Fired by webauthnp256-publickey-index.service's OnFailure= — the ONE moment
# the in-process Telegram alerting cannot speak for itself (the process is
# dead). Reads the same .env for the bot credentials; silently does nothing if
# Telegram is not configured.
[Unit]
Description=WebAuthn P256 index — service-failure Telegram alert

[Service]
Type=oneshot
EnvironmentFile=-/opt/webauthnp256-publickey-index/data/.env
# %H = hostname. Best-effort: 10s cap, failures ignored (alerting must never
# wedge the failure path itself).
ExecStart=/bin/sh -c 'test -n "$TELEGRAM_BOT_TOKEN" && test -n "$TELEGRAM_CHAT_ID" && curl -sf -m 10 -X POST "https://api.telegram.org/bot$TELEGRAM_BOT_TOKEN/sendMessage" -d chat_id="$TELEGRAM_CHAT_ID" --data-urlencode text="🛑 [webauthnp256-publickey-index] [%H] systemd unit FAILED (crashed/restart-limit). Creates are NOT being processed. journalctl -u webauthnp256-publickey-index -n 100" || true'
